# ops/supervisord.conf
[program:webapp]
command=uwsgi --socket 127.0.0.1:8080 --wsgi-file uwsgi.py --callable application  --processes 4  --threads 2 --logger syslog:uwsgi1234
directory=/home/deploy/blog
user=deploy